# profile.d/environment
export EDITOR=/usr/bin/vim
export http_proxy="http://172.27.10.114:3128"
export https_proxy="http://172.27.10.114:3128"
export ftp_proxy="http://172.27.10.114:3128"
